# Server/talent-relay.service
[Unit]
Description=Talent Trade Relay Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/relay
ExecStart=/usr/bin/python3 /root/relay/relay.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
